# install-files/etc/profile.d/lliurex-gsettings.sh
#!/bin/sh

# Disabled. Wait if is needed by KDE

#DOWNLOADS="$HOME/$(gettext xdg-user-dirs Downloads)"

#if [ ! -x "$DOWNLOADS" ]; then
#	DOWNLOADS=$HOME
#fi

#gsettings set org.blueman.transfer shared-path "$DOWNLOADS"

