# install-files/etc/profile.d/llx_sbin.sh
#!/bin/sh

# For now, let's make sure sbin is in path

export PATH=/usr/local/sbin:/usr/sbin:/sbin:/usr/games:$PATH



